# etc/systemd/system/docker.service.d/imagelocation.conf
[Service]
ExecStart= 
ExecStart=/usr/bin/docker daemon -g /home/docker/.data -H fd://
